# geoclue_fake.service
[Unit]
Description=Geoclue Service Faker

[Service]
Type=dbus
ExecStart=/usr/bin/geoclue_fake
BusName=org.freedesktop.GeoClue2
